# meta-ams/recipes-application/power-logger/power-logger/power-logger.service
[Unit]
Description=Solar Logger
After=network.target

[Service]
User=root
Group=root
ExecStart=/usr/bin/pcl --period 60 --storage /var/www/ams-ref-spa/html

[Install]
WantedBy=multi-user.target
